# crontab
0 6-23 * * *  /home/jollychang/tickets/check_tickets.py